# autostart/echo_sd-stop.path
[Unit]
Description=Watch for ECHO_SD unmount
After=local-fs.target

[Path]
PathChanged=/media/%u/
Unit=echo_sd-stop.service

[Install]
WantedBy=default.target

